# files/supervisord/oident.conf
[program:oidentd]
user=oident
command=/usr/sbin/oidentd -midf 
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
startretries=99999
