# files/supervisord/oident.conf
[program:oidentd]
command=/usr/sbin/oidentd -mf -u oident -g oident -id -a ::
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
startretries=99999
